# config/etc/services.d/cruise-control/run
#!/command/with-contenv sh

cd /opt/cruise-control
./kafka-cruise-control-start.sh config/cruisecontrol.properties 9091